# packages/base/all/rsyslog/after-install.sh
#!/bin/sh

systemctl enable rsyslog-config
