# packages/base/all/rsyslog/after-remove.sh
#!/bin/sh

systemctl disable rsyslog-config
