# packages/base/arm64/tai-dpll/before-remove.sh
#!/bin/sh

set -eux

systemctl disable tai-dpll.service
systemctl stop tai-dpll.service
