# packages/base/amd64/docker-dhcp-relay/after-install.sh
#!/bin/sh

set -eux

docker inspect docker-dhcp-relay:latest || docker load < /var/cache/apt/docker-dhcp-relay.gz
systemctl enable dhcp_relay
rm /var/cache/apt/docker-dhcp-relay.gz
