# packages/base/any/snmp/after-install.sh
#!/bin/sh

set -eux

systemctl enable gs-north-snmp.service
