# packages/base/amd64/docker-teamd/after-remove.sh
#!/bin/sh

docker inspect docker-teamd:latest || docker rmi docker-teamd:latest
systemctl disable teamd
