# packages/base/amd64/docker-fpm-frr/after-remove.sh
#!/bin/sh

docker inspect docker-fpm-frr:latest || docker rmi docker-fpm-frr:latest
systemctl disable bgp
