# packages/base/amd64/docker-database/after-remove.sh
#!/bin/sh

docker inspect docker-database:latest || docker rmi docker-database:latest
systemctl disable database
