# packages/base/any/tai/after-install.sh
#!/bin/sh

set -eux

systemctl enable tai.service
ldconfig -n /var/lib/tai/lib
