# packages/base/any/tai/before-remove.sh
#!/bin/sh

set -eux

systemctl disable tai.service
systemctl stop tai.service
