# packages/base/all/updategraph/after-remove.sh
#!/bin/sh

systemctl disable updategraph
