# packages/base/all/updategraph/after-install.sh
#!/bin/sh

systemctl enable updategraph
