# packages/base/amd64/usonic/builds/load-usonic-config.service
[Unit]
Description=load uSONiC platform config
After=k3s.service
Requires=k3s.service

[Service]
Type=oneshot
ExecStart=/usr/bin/load_usonic_config.sh
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
