# packages/base/amd64/ocnos/before-remove.sh
#!/bin/sh

set -eux

systemctl disable ocnos.service
ocnos.sh stop || true
