# packages/base/any/mgmt/builds/units/gs-mgmt-south.target
[Unit]
Description=Goldstone Management South Daemons

[Install]
WantedBy=gs-mgmt.target
