# packages/base/any/mgmt/builds/units/gs-mgmt-north.target
[Unit]
Description=Goldstone Management North Daemons

[Install]
WantedBy=gs-mgmt.target
