# packages/base/any/mgmt/builds/units/gs-mgmt-system.target
[Unit]
Description=Goldstone Management System Daemons

[Install]
WantedBy=gs-mgmt.target
