# packages/base/any/mgmt/builds/units/gs-mgmt-xlate.target
[Unit]
Description=Goldstone Management Translate Daemons

[Install]
WantedBy=gs-mgmt.target
