# packages/base/any/mgmt/builds/units/gs-south-ocnos-yang.service
[Unit]
Description=Goldstone Management South OcNOS Daemon YANG model setup
Before=gs-yang.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gs-yang.py --install south-ocnos
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
