# packages/base/amd64/docker-router-advertiser/after-remove.sh
#!/bin/sh

docker inspect docker-router-advertiser:latest || docker rmi docker-router-advertiser:latest
systemctl disable radv
