# packages/base/amd64/docker-router-advertiser/after-install.sh
#!/bin/sh

docker inspect docker-router-advertiser:latest || docker load < /var/cache/apt/docker-router-advertiser.gz
systemctl enable radv
rm /var/cache/apt/docker-router-advertiser.gz
