# packages/base/amd64/docker-router-advertiser/after-install.sh
#!/bin/sh

set -eux

IMAGE=docker-router-advertiser
docker inspect $IMAGE:latest || docker load < /var/cache/apt/$IMAGE.gz
systemctl enable radv
rm /var/cache/apt/$IMAGE.gz
