# packages/base/all/sonic/builds/pmon.service
[Unit]
Description=some description

[Service]
Type=oneshot
ExecStart=/bin/true