# packages/base/amd64/usonic/after-install.sh
#!/bin/sh

set -eux

systemctl daemon-reload
systemctl enable load-usonic-config.service

mkdir -p /var/lib/usonic/redis
