# packages/base/amd64/usonic/after-install.sh
#!/bin/sh

set -eux

systemctl enable usonic.service
