# packages/base/amd64/usonic/before-remove.sh
#!/bin/sh

set -eux

systemctl disable usonic.service
usonic.sh stop || true
