# packages/base/arm64/tai-gearbox/after-install.sh
#!/bin/sh

set -eux

systemctl enable tai-gearbox.service
