# packages/base/arm64/tai-gearbox/before-remove.sh
#!/bin/sh

set -eux

systemctl disable tai-gearbox.service
systemctl stop tai-gearbox.service
