# packages/base/amd64/docker-snmp-sv2/after-remove.sh
#!/bin/sh

docker inspect docker-snmp-sv2:latest || docker rmi docker-snmp-sv2:latest
systemctl disable snmp
