# packages/base/amd64/swss/after-install.sh
#!/bin/sh

systemctl enable swss
systemctl enable syncd
