# packages/base/amd64/swss/after-remove.sh
#!/bin/sh

systemctl disable swss
systemctl disable syncd
