# packages/base/amd64/docker-lldp-sv2/after-remove.sh
#!/bin/sh

docker inspect docker-lldp-sv2:latest || docker rmi docker-lldp-sv2:latest
systemctl disable lldp
