# packages/base/any/xlate-or/before-remove.sh
#!/bin/sh

set -eux

systemctl disable gs-xlate-or.service
systemctl stop gs-xlate-or.service
